# duckdns.service
[Unit]
Description=duckdns.py

[Service]
Type=simple
